# mystartupscript.service
[Unit] 
Description=Mystartupscript 
 
[Service] 
Type=oneshot
ExecStart=/YOUR/SCRIPT/GOES/HERE

[Install]
WantedBy=multi-user.target
